# build/edgebit-agent.service
[Unit]
Description=EdgeBit Agent

[Service]
ExecStart=/opt/edgebit/bin/edgebit-agent
EnvironmentFile=/etc/edgebit/config.env

[Install]
WantedBy=multi-user.target
