# build/edgebit-agent.service
[Unit]
Description=EdgeBit Agent

[Service]
ExecStart=/opt/edgebit/bin/edgebit-agent

[Install]
WantedBy=multi-user.target
