# dist/edgebit-agent.service
[Unit]
Description=EdgeBit Agent

[Service]
ExecStart=/opt/edgebit/edgebit-agent

[Install]
WantedBy=multi-user.target
