# dist/edgebit-agent.service
[Unit]
Description=EdgeBit Agent
Documentation=https://edgebit.io/docs/0.x/install-agent/
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/opt/edgebit/edgebit-agent
Restart=always

[Install]
WantedBy=multi-user.target
