# dist/edgebit-agent.service
[Unit]
Description=EdgeBit Agent
After=docker.service

[Service]
ExecStart=/opt/edgebit/edgebit-agent

[Install]
WantedBy=multi-user.target
